# etc/ansible/roles/twitter-contributors/templates/etc_supervisor_conf_d_twitter_contributors.conf.j2
[program:twitter_contributors]
command = node daemon.js
directory = {{ twitter_contributors_path }}